# ops/oracle/backflowverdict.service
[Unit]
Description=BackflowVerdict Spring Boot application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=backflow
Group=backflow
WorkingDirectory=/opt/backflowverdict
EnvironmentFile=/etc/backflowverdict/backflowverdict.env
ExecStart=/bin/bash -lc 'exec /usr/bin/java ${JAVA_OPTS:-} -jar /opt/backflowverdict/backflowverdict.jar'
SuccessExitStatus=143
Restart=always
RestartSec=5
TimeoutStopSec=20
StandardOutput=append:/opt/backflowverdict/logs/app.log
StandardError=append:/opt/backflowverdict/logs/app-error.log

[Install]
WantedBy=multi-user.target
